# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=/usr/local/bin/mysql-entrypoint.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mysql.err.log
stdout_logfile=/var/log/supervisor/mysql.out.log
user=root

[program:healthserver]
command=/usr/local/bin/health-server.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/health.err.log
stdout_logfile=/var/log/supervisor/health.out.log
user=root
